# deploy/tcc-bridge.service
[Unit]
Description=TCC Bridge - Telegram Claude Code Bridge
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/tcc-bridge
Environment=TCC_BRIDGE_CONFIG=/etc/tcc-bridge/bots.toml
EnvironmentFile=-/opt/tcc-bridge/.env
ExecStart=/opt/tcc-bridge/.venv/bin/python /opt/tcc-bridge/src/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tcc-bridge

[Install]
WantedBy=multi-user.target
