# deploy/tcc-bridge.service
[Unit]
Description=TCC Bridge - Telegram Claude Code Bridge
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/tcc-bridge
ExecStart=/usr/bin/python3 src/main.py
Restart=on-failure
RestartSec=10
EnvironmentFile=/home/ubuntu/tcc-bridge/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tcc-bridge

[Install]
WantedBy=multi-user.target
