# deploy/tcc-bridge.service
[Unit]
Description=TCC Bridge - Chat Claude Code Bridge
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/data/tcc/src
Environment=TCC_BRIDGE_CONFIG=/data/tcc/bots.toml
EnvironmentFile=-/data/tcc/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /data/tcc/src/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tcc-bridge

[Install]
WantedBy=multi-user.target
